# shr_backend/systemd/shr-backend.service
[Unit]
Description=SHR Camera Acquisition Backend
After=network.target usb-gadget.service
Wants=usb-gadget.service

[Service]
Type=simple
User=lvs
WorkingDirectory=/home/lvs/shr
ExecStart=/home/lvs/shr/SHR_Backend --simulator --gnss-port /dev/ttyUSB0 --output /home/lvs/frames
Restart=on-failure
RestartSec=3s
Environment="LD_LIBRARY_PATH=/home/lvs/vimba_libs"
Environment="GENICAM_GENTL64_PATH=/home/lvs/vimba_libs"
LimitRTPRIO=99
LimitMEMLOCK=infinity
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shr-backend

[Install]
WantedBy=multi-user.target
